# com.pi-sessions-viewer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pi-sessions-viewer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/setkyar/.pi/agent/bin/pi-sessions-viewer</string>
        <string>-p</string>
        <string>27183</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/pi-sessions-viewer.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/pi-sessions-viewer.error.log</string>

    <key>WorkingDirectory</key>
    <string>/tmp</string>
</dict>
</plist>
